# deploy/systemd/okx-trader.service
[Unit]
Description=OpenClaw OKX Trader v3.1
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/openclaw-trading-system/src
EnvironmentFile=-/root/openclaw-trading-system/.env
ExecStart=/usr/bin/python3 /root/openclaw-trading-system/src/ai_trader_v3_1.py --interval 60
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
