# deploy/systemd/site-ai-monitor.service
[Unit]
Description=OpenClaw Site AI Monitor (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/openclaw-trading-system/src
EnvironmentFile=-/root/openclaw-trading-system/.env
ExecStart=/usr/bin/python3 /root/openclaw-trading-system/src/site_platform_ai_monitor.py
